# ambulance.service
[Unit]
Description=LifeLine AI Flask App
After=network.target

[Service]
User=azureuser
WorkingDirectory=/home/azureuser/Lifeline-AI
Environment="PATH=/home/azureuser/Lifeline-AI/venv/bin"
EnvironmentFile=-/etc/lifeline-ai.env
ExecStart=/home/azureuser/Lifeline-AI/venv/bin/gunicorn --workers 2 --bind 0.0.0.0:5000 "app:create_app()"
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
